# measurement_buffer.service
[Unit]
Description=Measurement Buffer

[Service]
ExecStart=python3 /home/pi/measurement_buffer.py

[Install]
WantedBy=default.target
