# measurement_collector.service
[Unit]
Description=Measurement Collector

[Service]
ExecStart=python3 /home/pi/measurement_collector.py

[Install]
WantedBy=default.target
